# build/supervisord.conf
[unix_http_server]
; Socket im /tmp – beschreibbar durch den unprivilegierten 'claw'-User
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
; Kein 'user=' hier: supervisord läuft bereits als 'claw' (via Dockerfile USER)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; ---------------------------------------------------------------------------
; openclaw Gateway-Prozess
; Alle Container-Umgebungsvariablen (OPENCLAW_GATEWAY_PORT usw.) werden
; automatisch vererbt, da kein 'environment='-Override gesetzt ist.
; ---------------------------------------------------------------------------
[program:openclaw-gateway]
command=openclaw gateway --bind lan
autostart=true
autorestart=true
startretries=10
startsecs=5
; Logs direkt an Docker weiterleiten (docker logs <container>)
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

